# deb/dcgm-exporter-amd64/lib/systemd/system/dcgm-exporter.service
[Unit]
Description=DCGM Exporter
After=network.target

[Service]
ExecStart=/usr/bin/dcgm-exporter
KillMode=mixed
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
