# deb/otelcol-armhfv8/DEBIAN/postinst
#!/usr/bin/env bash

systemctl daemon-reload
systemctl enable otelcol
